# Server/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/entrypoint.sh mysqld

[program:flask]
command=python3 sql1.py
